# deploy/tencent-cloud/resume-assistant.service
[Unit]
Description=Resume Assistant Node Service
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/resume-assistant/server
Environment=NODE_ENV=production
Environment=HOST=127.0.0.1
EnvironmentFile=/opt/resume-assistant/server/.env
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=5
User=www-data
Group=www-data

[Install]
WantedBy=multi-user.target
